# systemd/orion-replica-sync.service
[Unit]
Description=Orion Sentinel External SSD Replication (Master -> Replica)
Documentation=https://github.com/orionsentinel/Orion-Sentinel-CoreSrv
After=local-fs.target
Wants=srv-orion-external_primary.mount srv-orion-external_replica.mount

[Service]
Type=oneshot
ExecStart=/opt/orion/Orion-Sentinel-CoreSrv/scripts/replicate-external.sh

# Run with low priority to minimize impact on other services
Nice=19
IOSchedulingClass=idle
IOSchedulingPriority=7

# Environment file for configuration
EnvironmentFile=-/opt/orion/Orion-Sentinel-CoreSrv/.env

# Timeouts for large syncs
TimeoutStartSec=86400

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orion-replica-sync

[Install]
WantedBy=multi-user.target
